# surveytrace-collector-ingest.service
[Unit]
Description=SurveyTrace collector ingest worker (master)
After=network.target
ConditionPathExists=/opt/surveytrace/data/surveytrace.db

[Service]
Type=simple
User=surveytrace
Group=surveytrace
SupplementaryGroups=www-data
WorkingDirectory=/opt/surveytrace/daemon
# Stable DB path (matches api/db.php) even if daemon/ is symlinked or script path differs.
Environment=SURVEYTRACE_INSTALL_DIR=/opt/surveytrace
ExecStart=/opt/surveytrace/venv/bin/python3 /opt/surveytrace/daemon/collector_ingest_worker.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/surveytrace/data

[Install]
WantedBy=multi-user.target
