# surveytrace-scheduler.service
[Unit]
Description=SurveyTrace scheduler daemon
After=network.target surveytrace-daemon.service
Documentation=https://github.com/yourorg/surveytrace
ConditionPathExists=/opt/surveytrace/data/surveytrace.db

[Service]
Type=simple
User=surveytrace
Group=surveytrace
SupplementaryGroups=www-data
WorkingDirectory=/opt/surveytrace/daemon
Environment=SURVEYTRACE_INSTALL_DIR=/opt/surveytrace
EnvironmentFile=-/etc/surveytrace/surveytrace.env
ExecStart=/opt/surveytrace/venv/bin/python3 scheduler_daemon.py
Restart=on-failure
RestartSec=10
# SQLite + subprocesses; default 1024 is easy to exhaust on long-lived schedulers.
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal
SyslogIdentifier=surveytrace-scheduler

# Same DB as scanner + WAL sidecars under data/; www-data shares group on SQLite files.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/surveytrace/data

[Install]
WantedBy=multi-user.target
