# surveytrace-credential-check-worker.service
[Unit]
Description=SurveyTrace credentialed check worker (SSH/SNMP execution; trusted-data observations)
After=network.target
ConditionPathExists=/opt/surveytrace/data/surveytrace.db

[Service]
Type=simple
User=surveytrace
Group=surveytrace
SupplementaryGroups=www-data
WorkingDirectory=/opt/surveytrace/daemon
Environment=SURVEYTRACE_INSTALL_DIR=/opt/surveytrace
EnvironmentFile=-/etc/surveytrace/surveytrace.env
ExecStart=/opt/surveytrace/venv/bin/python3 /opt/surveytrace/daemon/credential_check_worker.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/surveytrace/data

[Install]
WantedBy=multi-user.target
